# UNIQUE/aries/lxde/autostart/autostart-script.desktop
[Desktop Entry]
Type=Application
Name=Autostart Script
Comment=
Exec=sh -c '~/.config/autostart/resource/autostart.sh'

