# UNIQUE/aries/lxde/autostart/qbat.desktop
[Desktop Entry]
Type=Application
Name=qbat
Comment=
Exec=qbat
Terminal=false
